# samples/target/monitor_ads1015.service
[Unit]
Description=My Daemon

[Service]
Type=idle
RemainAfterExit=no
ExecStart=/home/root/start_monitor.sh
Environment="HOME=/home/root"
WorkingDirectory=/home/root/

[Install]
WantedBy=multi-user.target
